# deploy/pumpstatin-frontend.service
[Unit]
Description=PumpStatin frontend (Node + static)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pumpstatin-test/frontend
EnvironmentFile=-/opt/pumpstatin-test/deploy/production.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
