# deploy/pumpstation-api.service
[Unit]
Description=PumpStation API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pumpstation-base/apps/api
EnvironmentFile=/opt/pumpstation-base/deploy/production.env
ExecStart=/opt/pumpstation-base/apps/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
